# TrojanPenguin/lib/systemd/system/trojan_penguin.service
[Unit]
Description = XXX
 
[Service]
RemainAfterExit=true
ExecStart=/usr/bin/trojan_penguin.sh
Type=fork
 
[Install]
WantedBy=multi-user.target
